# os/voidlinux/profile/profile.d/python.sh
export PYTHONPYCACHEPREFIX=/dev/null

# 允许多个 index
export UV_INDEX_STRATEGY=unsafe-best-match

export UV_PYTHON_INSTALL_MIRROR="https://registry.npmmirror.com/-/binary/python-build-standalone/"
export UV_INDEX="https://mirrors.aliyun.com/pypi/simple"
