# os/voidlinux/profile/profile.d/python.sh
export PYTHONPYCACHEPREFIX=/dev/null

export UV_PYTHON_INSTALL_MIRROR="https://registry.npmmirror.com/-/binary/python-build-standalone/"
export UV_INDEX="https://mirrors.aliyun.com/pypi/simple"