# os/freebsd/profile.d/lang.sh
# gui 的脚本会自动修改，这里保持注释
#export LANG=zh_CN.UTF-8
