# os/freebsd/profile.d/ls.sh
alias ll='ls -la -G' 2>/dev/null
alias l.='ls -G -d .*' 2>/dev/null
alias ls='ls -G' 2>/dev/null
