# os/freebsd/profile.d/ls.sh
alias ls='ls -G' 2>/dev/null
alias ll='ls -l -G' 2>/dev/null
