# os/freebsd/profile.d/nodejs.sh
export MISE_NODE_MIRROR_URL=https://npmmirror.com/mirrors/node/
