# os/freebsd/profile.d/proxy.sh
alias pon="export http_proxy=http://127.0.0.1:8118 https_proxy=http://127.0.0.1:8118 all_proxy=socks5://127.0.0.1:1080"
alias poff="unset http_proxy https_proxy all_proxy"

alias proxy="proxychains"