# os/freebsd/zsh/user/etc/profile.d/tls.sh
export SSLKEYLOGFILE=/tmp/ssl_key_$USER.log
