# os/debian/profile/profile.d/erlang.sh
export ERL_AFLAGS="-kernel shell_history enabled"
export HEX_CDN=https://repo.hex.pm
