# os/debian/profile/profile.d/http.sh
export HTTP_CONNECT_TIMEOUT=36000
export HTTP_RECEIVE_TIMEOUT=36000