# os/manjaro/profile.d/elixir.sh
export HEX_UNSAFE_HTTPS=1
export HEX_MIRROR="https://repo.hex.pm"
