# os/manjaro/profile.d/pacman-mirror.sh
alias pacman-mirror-select="sudo pacman-mirrors -i --geoip -m rank"
